# pi/systemd/berry-frontend.service
[Unit]
Description=Berry Frontend
After=network-online.target berry-backend.service
Wants=network-online.target
Requires=berry-backend.service

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/berry/frontend
ExecStart=/usr/bin/npx serve -s dist -l 3000
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=HOME=/home/admin

[Install]
WantedBy=multi-user.target

